# drawreportDeploy/drawreport-worker.service
# /etc/systemd/system/drawreport-worker.service
# Background report worker (paid -> generating -> delivered). One instance.
[Unit]
Description=DrawReport report worker
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/DrawReport
ExecStart=/var/www/DrawReport/venv/bin/python worker.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
